# deploy/chpc-gpu-finder.service
[Unit]
Description=CHPC GPU Finder web server
After=network.target

[Service]
Type=simple
WorkingDirectory=%h/chpc-gpu-finder
ExecStart=/usr/bin/env python3 %h/chpc-gpu-finder/server.py --host 0.0.0.0 --port 8000
EnvironmentFile=%h/chpc-gpu-finder/.env
Restart=always
RestartSec=5

[Install]
WantedBy=default.target
